# conf/supervisor.conf
[program:govtrack_web]
command = /home/govtrack/web/conf/gunicorn_start www 20
directory = /home/govtrack/web
user = govtrack
env = LANG=en_US.UTF-8
stopsignal = QUIT
