# conf/supervisor.ini
[program:www_govtrack_us]
command = /home/govtrack/web/.env/bin/uwsgi --wsgi-file wsgi.py --socket /tmp/uwsgi_govtrack_www.sock --chmod-socket=666 --pidfile /tmp/uwsgi_govtrack_www.pid --processes 12 --cheaper 2
directory = /home/govtrack/web
user = govtrack

[program:test_govtrack_us]
command = /home/govtrack/web_dev/.env/bin/uwsgi --wsgi-file wsgi.py --socket /tmp/uwsgi_govtrack_test.sock --chmod-socket=666
directory = /home/govtrack/web_dev
user = govtrack

