# conf/supervisor.conf
[program:govtrack_web]
command = /home/govtrack/web/conf/uwsgi_start www 10
directory = /home/govtrack/web
user = govtrack
env = LANG=en_US.UTF-8

#[program:test_govtrack_us]
#command = /home/govtrack/web_dev/conf/uwsgi_start test 2
#directory = /home/govtrack/web_dev
#user = govtrack
#env = LANG=en_US.UTF-8
